# kiro-devops.service
[Unit]
Description=Feishu Kiro Bot Bridge (WebSocket)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/kiro-devops
EnvironmentFile=/home/ubuntu/kiro-devops/.env
Environment=PATH=/home/ubuntu/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/python3 /home/ubuntu/kiro-devops/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
